# msgme.socket
[Socket]
ListenStream=27087

[Install]
WantedBy=default.target
